# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:gunicorn]
command=/usr/local/bin/gunicorn --bind 127.0.0.1:8000 --workers 3 healthstack.wsgi:application
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/gunicorn.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx.log